# com.example.calendar.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.calendar.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Path/to/Calendar/env/bin/python</string>
    <string>/Path/to/Calendar/main.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Path/to/Calendar</string>

  <key>StartCalendarInterval</key>
  <array>
    <!-- Normal time (07:30) on all days except Wednesday and weekends, this is my personal configuration -->
    <dict><key>Weekday</key><integer>0</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>

    <!-- Wednesday and weekends at 09:00 -->
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>

  </array>

  <key>StandardOutPath</key>
  <string>/Path/to/Calendar/calendar.log</string>
  <key>StandardErrorPath</key>
  <string>/Path/to/Calendar/calendar.err.log</string>
</dict>
</plist>